# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[program:api]
command=/root/.openharness-venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/opt/oh-service
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api-stdout.log
stderr_logfile=/var/log/supervisor/api-stderr.log
environment=PYTHONPATH="/app/src:/opt/oh-service"

[program:worker]
command=/root/.openharness-venv/bin/celery -A app.workers.celery_app.celery_app worker -l info -c 4
directory=/opt/oh-service
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/worker-stdout.log
stderr_logfile=/var/log/supervisor/worker-stderr.log
environment=PYTHONPATH="/app/src:/opt/oh-service"

[program:beat]
command=/root/.openharness-venv/bin/celery -A app.workers.celery_app.celery_app beat -l info
directory=/opt/oh-service
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/beat-stdout.log
stderr_logfile=/var/log/supervisor/beat-stderr.log
environment=PYTHONPATH="/app/src:/opt/oh-service"

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface
